# deploy/ansible/roles/web_certs/templates/certbot.timer.j2
[Unit]
Description=Renew certificates for {{ app_name }} every day

[Timer]
OnCalendar=daily
RandomizedDelaySec=1day
Persistent=true
